# dot_config/autostart/startup.desktop
[Desktop Entry]
Type=Application
Exec=startup
Name=My Startup Applications
