# dot_config/autostart/terminal-start-server.desktop
[Desktop Entry]
Type=Application
Exec=terminal start-server
Name=My Terminal Server
