# dot_config/autostart/wallpaper.desktop
[Desktop Entry]
Type=Application
Exec=change-wallpaper
Name=Change wallpaper every 15 mins
